# cookbooks/kings/templates/default/kings.supervisor.conf.erb
[program:kings]
command=<%= @kings_bin %> -c <%= @kings_config %> 
autostart=true
autorestart=unexpected
user=<%= @kings_user %>
redirect_stderr=true
stdout_logfile=<%= @kings_log %>


